# systemd/housekeeper.timer
[Unit]
Description=Run laptop health check daily at 06:00 IST

[Timer]
OnCalendar=*-*-* 06:00:00
# Persistent=true: if 06:00 was missed (machine asleep/off), run once at
# next wake/boot. Won't wake a suspended machine (no WakeSystem).
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
